# systemd/traid-worker.service
[Unit]
Description=TRAID Privileged Worker Daemon
Documentation=https://github.com/your-org/traid
After=network.target

[Service]
Type=simple
ExecStart=/usr/lib/traid/venv/bin/python -m worker.daemon
WorkingDirectory=/usr/lib/traid
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=traid-worker

# Security hardening
PrivateTmp=yes
NoNewPrivileges=no

[Install]
WantedBy=multi-user.target
