# systemd/traid-web.service
[Unit]
Description=TRAID Web API (Uvicorn/FastAPI)
Documentation=https://github.com/albal/traid
After=network.target traid-worker.service
Requires=traid-worker.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/usr/lib/traid
ExecStart=/usr/lib/traid/venv/bin/uvicorn api.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=traid-web

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/usr/lib/traid

[Install]
WantedBy=multi-user.target
