# debian/postinst
#!/bin/sh
set -e

TRAID_DIR=/usr/lib/traid

case "$1" in
    configure|reconfigure)
        # (Re)create the Python venv. Safe to run on upgrade — venv is
        # regenerated from scratch so stale bytecode is never an issue.
        python3.13 -m venv "$TRAID_DIR/venv"

        # btrfsutil is a system C extension (python3-btrfsutil) not on PyPI.
        # Copy it into the venv so the worker daemon can import it as root.
        BTRFSUTIL_SO=$(find /usr/lib/python3 -name 'btrfsutil.cpython-*.so' 2>/dev/null | head -1)
        if [ -n "$BTRFSUTIL_SO" ]; then
            cp "$BTRFSUTIL_SO" "$TRAID_DIR/venv/lib/python3.13/site-packages/"
        fi

        # Install traid and all dependencies from the pre-built wheels bundled in the package.
        "$TRAID_DIR/venv/bin/pip" install \
            --no-index \
            --find-links "$TRAID_DIR/vendor" \
            traid \
            --quiet

        # Create an upload spool directory on the same filesystem as the
        # ISO destination, writable by the web user. Uploaded ISOs are
        # streamed here (not /tmp, which is tmpfs and too small) and the
        # worker then does an atomic rename into /var/lib/traid/iso/.
        mkdir -p /var/lib/traid/iso-upload
        chown www-data:www-data /var/lib/traid/iso-upload
        chmod 0750 /var/lib/traid/iso-upload

        # Point lighttpd's document root at the TRAID static files.
        # lighty-enable-mod only handles numbered module configs, so we
        # manage the symlink and document-root ourselves.
        if [ -f /etc/lighttpd/lighttpd.conf ]; then
            sed -i 's|server\.document-root\s*=\s*"/var/www/html"|server.document-root = "/usr/share/traid"|' \
                /etc/lighttpd/lighttpd.conf
        fi

        # Enable the TRAID lighttpd configuration.
        ln -sf /etc/lighttpd/conf-available/traid.conf \
               /etc/lighttpd/conf-enabled/traid.conf

        if [ -d /run/systemd/system ]; then
            systemctl reload-or-restart lighttpd.service || true
        fi

        # Reload unit files (handles both fresh install and upgrade).
        if [ -d /run/systemd/system ]; then
            systemctl daemon-reload

            systemctl enable traid-worker.service traid-web.service

            # Restart if already running (upgrade), otherwise start fresh.
            if systemctl is-active --quiet traid-worker.service; then
                systemctl restart traid-worker.service
                systemctl restart traid-web.service
            else
                systemctl start traid-worker.service
                # Give the worker a moment to bind the socket.
                sleep 1
                systemctl start traid-web.service || true
            fi
        fi
        ;;
esac

#DEBHELPER#
exit 0
