# debian/postrm
#!/bin/sh
set -e

case "$1" in
    remove)
        if [ -d /run/systemd/system ]; then
            systemctl disable traid-worker.service traid-web.service || true
            systemctl daemon-reload || true
        fi
        ;;
    purge)
        # Remove generated venv (source files are handled by dpkg).
        rm -rf /usr/lib/traid/venv
        if [ -d /run/systemd/system ]; then
            systemctl daemon-reload || true
        fi
        ;;
esac

#DEBHELPER#
exit 0
